# dockerfiles/theia/supervisord.conf
[supervisord]

[program:theia]
command=node /theia_launcher/theia_launcher.js
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
