# dockerfiles/theia/supervisord.conf
[supervisord]

[program:theia]
command=node /theia_launcher/theia_launcher.js
stopasgroup=true
killasgroup=true
